# launchd/com.axinova.agent-tech-writer.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.axinova.agent-tech-writer</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/axinova-agent/workspace/axinova-agent-fleet/scripts/agent-launcher.sh</string>
        <string>tech-writer</string>
        <string>/Users/axinova-agent/workspace</string>
        <string>180</string>
    </array>
    <key>UserName</key>
    <string>axinova-agent</string>
    <key>WorkingDirectory</key>
    <string>/Users/axinova-agent/workspace</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>HOME</key>
        <string>/Users/axinova-agent</string>
    </dict>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>
    <key>ThrottleInterval</key>
    <integer>60</integer>
    <key>StandardOutPath</key>
    <string>/Users/axinova-agent/logs/agent-tech-writer-stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/axinova-agent/logs/agent-tech-writer-stderr.log</string>
</dict>
</plist>
